# iso/rootfs-files/etc/systemd/system/run-besconf.automount
[Unit]
Description=BES Configuration Partition Automount
ConditionPathExists=/dev/disk/by-partuuid/e2bac42b-03a7-5048-b8f5-3f6d22100e77

[Automount]
Where=/run/besconf
TimeoutIdleSec=60

[Install]
WantedBy=local-fs.target
