# iso/rootfs-files/etc/systemd/system/run-besconf.automount
[Unit]
Description=BES Configuration Partition Automount
ConditionPathExists=/dev/disk/by-label/BESCONF

[Automount]
Where=/run/besconf
TimeoutIdleSec=60

[Install]
WantedBy=local-fs.target
