# image/files/systemd/setup-tpm-unlock.service
[Unit]
Description=Setup TPM2 auto-unlock for LUKS
After=local-fs.target luks-reencrypt.service
ConditionPathExists=/dev/tpm0
ConditionPathExists=/dev/tpmrm0
ConditionPathExists=!/etc/luks/tpm-enrolled

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-tpm-unlock
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
